# Iot Code (DO NOT TOUCH)/hls-uploader.service
[Unit]
Description=HLS File Uploader - Uploads HLS segments to VPS
After=network-online.target ffmpeg-hls.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/raspi-device-manager
ExecStart=/usr/bin/python3 /home/pi/raspi-device-manager/hls_uploader.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

